# ibgateway/ibgateway.service
[Unit]
Description=IB Gateway 24/7 (with VNC for noVNC login)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:99

# Start script: Xvfb + x11vnc + websockify + IB Gateway
ExecStart=/bin/bash /root/wall-clock/ibgateway/start_ibgateway.sh

# Keep running 24/7: restart on crash or exit
Restart=always
RestartSec=30

StandardOutput=append:/var/log/ibgateway.log
StandardError=append:/var/log/ibgateway.log

[Install]
WantedBy=multi-user.target
